# home/dot_config/systemd/user/waybar.service.d/override.conf
[Service]
Type=exec
ExecCondition=/lib/systemd/systemd-xdg-autostart-condition "wlroots:sway:Wayfire:labwc:Hyprland" ""
Slice=app-graphical.slice
