# home/dot_config/systemd/user/vvn_hyprstart.target
[Unit]
Description=Hyprland DE dependencies
Documentation=man:systemd.special(7)
RefuseManualStart=no
StopWhenUnneeded=no
After=dms.service \
      elephant.service \
      walker.service \
      hyprpolkitagent.service \
      hyprpaper.service \
      hypridle.service \
      dropbox.service \
      blueman-applet.service
Wants=dms.service \
      elephant.service \
      walker.service \
      hyprpolkitagent.service \
      hyprpaper.service \
      hypridle.service \
      dropbox.service \
      blueman-applet.service
